# files/rkhunter.cron
0 4 * * * /usr/bin/rkhunter -c --cronjob